# minio/rootfs/etc/services.d/minio/run
#!/command/with-contenv bash
set -e

echo "[minio] Starting MinIO server..."

exec /usr/bin/minio server /share/minio \
    --address :9000 \
    --console-address :9002
