# wallos/rootfs/etc/services.d/wallos-nginx/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting Wallos backend nginx (port 8080)..."

exec nginx -c /etc/wallos-nginx.conf -g 'daemon off;'
