# wallos/rootfs/etc/services.d/wallos-nginx/run
#!/command/with-contenv bash
set -e

echo "[wallos-nginx] Starting Wallos backend nginx (port 8080)..."

exec nginx -c /etc/wallos-nginx.conf -g 'daemon off;'
