# wallos/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting crond..."

exec crond -f
