# wallos/rootfs/etc/services.d/crond/run
#!/command/with-contenv bash
set -e

echo "[crond] Starting crond..."

exec crond -f -L /dev/stdout
