# wallos/rootfs/etc/services.d/nginx-proxy/run
#!/command/with-contenv bash
set -e

echo "[nginx-proxy] Waiting for Wallos backend on port 8080..."

# Simple wait loop
for i in {1..60}; do
    if nc -z localhost 8080 2>/dev/null; then
        echo "[nginx-proxy] Wallos backend is ready!"
        break
    fi
    sleep 1
done

echo "[nginx-proxy] Starting nginx reverse proxy..."

exec nginx
