# wallos/rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv bash
set -e

echo "[php-fpm] Starting PHP-FPM..."

exec php-fpm -F
